# config/systemd/bsc-disk-check.service
[Unit]
Description=BSC node /data disk-usage check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bsc-disk-check.sh
